# sources/home/mk/.config/autostart/lazycast.desktop
[Desktop Entry]
Name=lazycast
Exec=/opt/lazycast/allnew.sh
Type=application 
Terminal=true
